# metadata/PhoenixMaster.conf
start on startup

expect fork

script
exec sudo /home/ubuntu/bajabeaglebone/startup.sh &
end script
