# metadata/PhoenixMaster.conf
start on runlevel [2345]
stop on runlevel [^2345]

description "PhoenixMaster"

script
chdir /home/ubuntu/bajabeaglebone
exec sudo python PhoenixMaster.py
end script
